# airalertd.service
[Unit]
Description=Air alert notification service

[Service]
ExecStart=/usr/local/bin/airalert/airalert
User=root

[Install]
WantedBy=default.target
